# src/main/assembly/dist/install/dd-ingest-flow.service
[Unit]
Description=Dd Ingest Flow Service
After=payara.service easy-validate-dans-bag.service

[Service]
ExecStart=/opt/dans.knaw.nl/dd-ingest-flow/bin/dd-ingest-flow server /etc/opt/dans.knaw.nl/dd-ingest-flow/config.yml
# Java returns 143 even if the SIGTERM was handled correctly.
SuccessExitStatus=143

User=dd-ingest-flow
Group=dd-ingest-flow

[Install]
WantedBy=multi-user.target
